# files/etc/rc.d/xray
#!/bin/sh

# PROVIDE: xray
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="xray"
rcvar="xray_enable"
command="/usr/local/bin/xray"
config_file="/usr/local/etc/xray/config.json"
command_args="-config ${config_file}"
pidfile="/var/run/${name}.pid"
required_files="${config_file}"
log_file="/var/log/xray_service.log"  # Optional log file
load_rc_config $name
: ${xray_enable:="NO"}

start_precmd="xray_precmd"

xray_precmd() {
    if [ ! -d "/var/run" ]; then
        mkdir -p /var/run
    fi
    if [ ! -d "/var/log" ]; then
        mkdir -p /var/log
    fi
}

start_cmd="${command} ${command_args} -log ${log_file} -pidfile ${pidfile}"

run_rc_command "$1"
